# supervisor/conf.d/fluxbox.conf
[program:fluxbox]
command=/usr/bin/fluxbox
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
stopsignal=KILL
