# app/grlevelx.conf
[program:grlevelx]
command=/opt/grlevelx.sh
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
stopsignal=KILL